# systemd/usbipd.service
[Unit]
Description=USB/IP Daemon
After=network.target

[Service]
ExecStart=/usr/sbin/usbipd -D
Restart=always

[Install]
WantedBy=multi-user.target
